# ops/cron/vaultfire_reinforcement.cron
# Regenerate Vaultfire SLA and service map every Friday at 04:04 UTC
4 4 * * 5 /usr/bin/env python -m vaultfire.dependencies.refresh
